# overlay/etc/cron.d/snipe-it
# Snipe IT cron job to run any scheduled tasks every minute
* * * * * root /usr/local/bin/turnkey-artisan schedule:run >> /dev/null 2>&1
